# services/com.titan.ideas.drain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.titan.ideas.drain</string>

  <key>ProgramArguments</key>
  <array>
    <string>__HARNESS_PATH__/bin/idea-drain.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__SESSION_DIR__/ideas-drain.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__SESSION_DIR__/ideas-drain.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>LC_ALL</key>
    <string>en_US.UTF-8</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
