# services/io.amg.credential-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  io.amg.credential-backup — Mac-side LaunchAgent
  DR-AMG-SECURITY-01 Phase 2 Task 2.3 companion: nightly GPG+R2 credential
  inventory backup from Mac (separate from VPS amg-backup.sh which covers
  app state only, not operator credential inventory).

  Schedule: daily 23:00 local time.
  Install: cp services/io.amg.credential-backup.plist ~/Library/LaunchAgents/
           launchctl load -w ~/Library/LaunchAgents/io.amg.credential-backup.plist
  Verify:  launchctl list | grep io.amg.credential-backup

  Lockout-risk gate: all 5 = NO. Reads local credential inventory, GPG-encrypts,
  uploads to separate R2 bucket. No SSH/firewall/auth write surface.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.amg.credential-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/solonzafiropoulos1/titan-harness/bin/credential-backup-r2.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/solonzafiropoulos1/Library/Logs/amg/credential-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/solonzafiropoulos1/Library/Logs/amg/credential-backup.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/solonzafiropoulos1/titan-harness</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
