# services/titan-ideas-drain.timer
[Unit]
Description=Run Titan Ideas Drainer every 60s

[Timer]
OnBootSec=45
OnUnitActiveSec=60
AccuracySec=5
Persistent=true
Unit=titan-ideas-drain.service

[Install]
WantedBy=timers.target
