# opt-amg-titanium/systemd/amg-titanium-post-mortem.timer
[Unit]
Description=Run amg-titanium-post-mortem daily at 04:30 UTC

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
Persistent=true
Unit=amg-titanium-post-mortem.service

[Install]
WantedBy=timers.target
