# launchd/com.amg.hercules-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amg.hercules-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/solonzafiropoulos1/titan-harness/scripts/hercules_daemon.py</string>
        <string>--watch</string>
        <string>--interval</string>
        <string>30</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/solonzafiropoulos1/titan-harness</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/hercules_daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/hercules_daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/solonzafiropoulos1</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
