# launchd/com.amg.warden-enforcer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amg.warden-enforcer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/solonzafiropoulos1/titan-harness/scripts/warden_enforcer.py</string>
        <string>--watch</string>
        <string>--interval</string>
        <string>600</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MCP_BASE</key>
        <string>https://memory.aimarketinggenius.io</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/warden_enforcer.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/warden_enforcer.err.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/solonzafiropoulos1/titan-harness</string>
</dict>
</plist>
