# launchd/com.amg.hercules-onboarding.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amg.hercules-onboarding</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/solonzafiropoulos1/titan-harness/scripts/hercules_onboarding_paste.py</string>
        <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/solonzafiropoulos1/titan-harness</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/hercules_onboarding.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/solonzafiropoulos1/.openclaw/logs/hercules_onboarding.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/solonzafiropoulos1</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
