# launchd/io.titan.nosleep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.titan.nosleep</string>

    <!-- caffeinate -s prevents system sleep on AC power; display can still
         sleep but SSH remains alive. CT-0415-05 Phase 3. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Logs in case caffeinate exits unexpectedly. -->
    <key>StandardOutPath</key>
    <string>__HOME__/.titan-nosleep.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.titan-nosleep.stderr.log</string>

    <!-- Low priority. -->
    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
